# deploy/systemd/ocr-hvks-api.service
[Unit]
Description=OCR_HVKS combined FastAPI (OCR + Extract + LaTeX)
After=network-online.target vllm@%i.service
Wants=network-online.target vllm@%i.service

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/ocr_hvks
EnvironmentFile=-/opt/ocr_hvks/.env
ExecStart=/usr/bin/env bash /opt/ocr_hvks/deploy/start_api.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
