# deploy/systemd/ocr-hvks-api.service
[Unit]
Description=OCR_HVKS combined FastAPI (OCR + Extract + LaTeX)
After=network-online.target vllm@%i.service
Wants=network-online.target vllm@%i.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/AI_project/OCR_HVKS
EnvironmentFile=-/home/%i/AI_project/OCR_HVKS/.env
ExecStart=/usr/bin/env bash /home/%i/AI_project/OCR_HVKS/deploy/start_api.sh
Restart=on-failure
RestartSec=5
# Chờ uvicorn graceful shutdown (xem API_TIMEOUT_GRACEFUL) rồi mới SIGKILL.
TimeoutStopSec=40
# Trần file descriptor: connection pool (LLM_CONCURRENCY) + nhiều SSE đồng
# thời dễ vượt mặc định 1024 → cần nới rộng.
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
