# deploy/systemd/vllm.service
[Unit]
Description=vLLM server (Qwen3.6) for OCR_HVKS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/ocr_hvks
EnvironmentFile=-/opt/ocr_hvks/.env
ExecStart=/usr/bin/env bash /opt/ocr_hvks/deploy/start_vllm.sh
Restart=on-failure
RestartSec=10
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
